# 42-cursus/common-core-rank01/Born2beroot/src/monitoring.timer
[Unit]
Description='monitoring' broadcast service timer

[Timer]
OnBootSec=10min
OnUnitActive=10min
AccuracySec=1msec

[Install]
WantedBy=timers.target
