# 42-cursus/common-core-rank01/Born2beroot/src/monitoring.service
[Unit]
Description=Run the 'monitoring' broadcast script

[Service]
ExecStart=monitoring.sh
